# jenkins.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type="manifest" name="Jenkins">
    <service name="application/jenkins" type="service" version="1">
        <create_default_instance enabled="false"/>
        <single_instance/>
        
        <method_context>
            <method_credential user='jenkins' group='jenkins'/>
            <method_environment>
                <envvar name='PATH' value='/usr/bin:/usr/local/bin'/>
            </method_environment>
        </method_context>

        <exec_method type="method" name="start" exec="/opt/jenkins/start_jenkins.sh" timeout_seconds="0"/>
        <exec_method type="method" name="stop" exec=":kill -TERM" timeout_seconds="30"/>

        <!-- We are going to be kicking off a single child process so we want Wait mode-->
        <property_group name='startd' type='framework'>
            <propval name='duration' type='astring' value='child'/>
        </property_group>

        <stability value="Evolving"/>

        <template>
            <common_name>
                <loctext xml:lang='C'>Jenkins Continuous Build Server</loctext>
            </common_name>
            <documentation>
                <doc_link name='jenkins-ci.org' uri='http://jenkins-ci.org/'/>
            </documentation>
        </template>
    </service>
</service_bundle>
